# syndicator/crontab
* * * * *  node ~/Documents/classes/programming/pulsd/syndicator/dst/