# crontab
* * * * *  node ~/Documents/classes/programming/pulsd/syndicator/dst/ 2>&1 ~/Documents/classes/programming/pulsd/cronlog