# crontab
*/5 * * * * /usr/local/bin/node ~/pulsd/syndicator/dst/index.js